# rc.local
#!/bin/bash

# secure the system
ufw allow 22
ufw allow 80
ufw enable

# phone home
/usr/bin/curl http://IPADDR:8080/boot?name=NAME

# catalogue the filesystem
/usr/bin/updatedb

# schedule for phoning home
cat > /tmp/rootcron <<EOF
* * * * * /usr/bin/curl http://IPADDR:8080/cron?name=NAME
EOF
crontab -u vm /tmp/rootcron
rm /tmp/rootcron

# run boot script
cat > /tmp/boot.sh <<EOF
RUN
EOF
chmod u+x /tmp/boot.sh
/tmp/boot.sh
rm /tmp/boot.sh

exit 0
