# rc.local
#!/bin/bash -e
echo "howdy" > /root/test.txt
exit 0
