# rc.local
#!/bin/bash -e
ufw allow 22
ufw allow 80
ufw enable
/usr/bin/updatedb
cat > /tmp/rootcron <<EOF
* * * * * /usr/bin/curl http://192.168.1.2:8080/cron?name=NAME
EOF
crontab -u vm /tmp/rootcron
rm /tmp/rootcron
/usr/bin/curl http://192.168.1.2:8080/boot?name=NAME
exit 0
